# supervisor/conf.d/lxpanel.conf
[program:lxpanel]
priority=15
directory=%HOME%
command=/usr/bin/lxpanel --profile LXDE
user=%USER%
environment=DISPLAY="%DISPLAY%",HOME="%HOME%",USER="%USER%"
autorestart=true
