# supervisor/conf.d/pcmanfm.conf
[program:pcmanfm]
priority=15
directory=%HOME%
command=/usr/bin/pcmanfm --desktop --profile LXDE
user=%USER%
environment=DISPLAY="%DISPLAY%",HOME="%HOME%",USER="%USER%"
autorestart=false
